# serverconf/production/etc/init/sh8recv.conf
description "sh8email mail receveing server"

start on runlevel [2345]
stop on runlevel [06]

respawn

script
	export SH8EMAIL_DJANGO_SECRET_KEY={{SH8EMAIL_DJANGO_SECRET_KEY}}
	export DJANGO_SETTINGS_MODULE=sh8email.settings_prod
	export SH8EMAIL_POSTGRES_PW={{SH8EMAIL_POSTGRES_PW}}

	exec /home/sh8email/venv/bin/python3 /home/sh8email/sh8email-django/manage.py runrecv >> /var/log/sh8recv.log 2>> /var/log/sh8recv.err
end script
